# scav.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/sbin/screenscav.sh

[Install]
WantedBy=default.target
